# webpasswd.service
[Unit]
Description=webpasswd — web-based password changer
Documentation=https://github.com/csett86/webpasswd
After=network.target

[Service]
Type=simple

# pam_unix operates in direct mode (no setuid helpers) only when UID=0.
# unix_update, the PAM helper for writing /etc/shadow, carries no setuid bit
# and therefore also requires root-level file access from its caller.
User=root
Group=root

# Minimal Linux capabilities required for PAM to authenticate users and update
# /etc/shadow directly:
#
#   CAP_DAC_READ_SEARCH  — open /etc/shadow for reading (bypasses DAC read check)
#   CAP_DAC_OVERRIDE     — open /etc/shadow for writing and create the PAM lock
#                          file /etc/.pwd.lock in the root-owned /etc/ directory
#   CAP_AUDIT_WRITE      — record password-change events in the kernel audit log
#
# All other root capabilities (CAP_SYS_ADMIN, CAP_NET_ADMIN, CAP_SYS_PTRACE,
# CAP_SETUID, CAP_SETGID, …) are stripped by the bounding set below.
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE CAP_AUDIT_WRITE
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE CAP_AUDIT_WRITE

# Prevent the process from gaining additional privileges via setuid/setgid
# binaries. Safe because pam_unix uses the direct code path (UID=0), so the
# setuid helper binaries are never invoked.
NoNewPrivileges=true

# Path to the binary.
ExecStart=/usr/local/bin/webpasswd \
    -addr :8080 \
    -rate-limit 5 \
    -rate-window 15m \
    -x-forwarded-for true

# Working directory must contain templates/ and static/ subdirectories
# (or adjust the binary to accept explicit paths).
WorkingDirectory=/etc/webpasswd

Restart=on-failure
RestartSec=5

# Lock down the filesystem; only the shadow file and its PAM lock/backup need
# to be writable at runtime.
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/shadow /etc/shadow- /etc/.pwd.lock
PrivateTmp=true

# Log to the journal.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webpasswd

[Install]
WantedBy=multi-user.target
